# Mangdang/System/rc.local
#!/bin/bash
#
# 1.Initialize PWM channels
#
# Copyright (c) 2020, Mangdang Technology Co., Limited
#
#

# export pwm channel
for i in $(seq 0 15)
do
	echo $i > /sys/class/pwm/pwmchip0/export
	echo  4000000 > /sys/class/pwm/pwmchip0/pwm$i/period
done

# set cpu performance
echo performance > /sys/bus/cpu/devices/cpu0/cpufreq/scaling_governor

mpg123 /home/ubuntu/Music/power_on.mp3 &
